# playbooks/files/iris-stage.service
[Unit]
Description=Staging client for iris snap deployment pipeline

[Service]
Type=oneshot
ExecStart=/home/ubuntu/.local/bin/stage

[Install]
WantedBy=multi-user.target